# tools/ci/systemd/nebula-cleanup.timer
[Unit]
Description=Daily cleanup timer for Nebula builds

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
