# tools/ci/systemd/nebula-cleanup.service
[Unit]
Description=Nebula cleanup job

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nebula-cleanup-builds.sh
