# tools/ci/systemd/nebula-runner-update.service
[Unit]
Description=Nebula runner update job

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nebula-runner-update.sh
